# docker/supervisor/vozor-laravel-worker.conf
[program:voroz-laravel-worker]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www/voroz.serv
command='php artisan queue:work sync --sleep=3 --tries=3 --max-time=3600'
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=root
numprocs=5
redirect_stderr=true
stdout_logfile=/var/www/voroz.serv/logs/worker.log
stopwaitsecs=3600
